# quadlet/dify-nginx.container
[Unit]
Description=Dify Nginx Reverse Proxy
After=network-online.target

[Container]
Image=docker.io/nginx:alpine
ContainerName=dify-nginx
Pod=dify-stack.pod

Volume=%h/.local/share/ai-models/dify/nginx.conf:/etc/nginx/nginx.conf:Z,ro

# -s without -f: exit 0 on any HTTP response (including 502 while upstreams
# are still starting), exit 1 only if nginx itself is not responding
HealthCmd=curl -s -o /dev/null http://localhost/ || exit 1
HealthInterval=30s
HealthStartPeriod=30s

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=60

[Install]
WantedBy=default.target
