# quadlet/dify-nginx.container
[Unit]
Description=Dify Nginx Reverse Proxy
After=network-online.target

[Container]
Image=docker.io/nginx:alpine
ContainerName=dify-nginx
Pod=dify-stack.pod

Volume=%h/.local/share/ai-models/dify/nginx.conf:/etc/nginx/nginx.conf:Z,ro

# nginx:alpine has wget (busybox) but not curl
HealthCmd=wget -q --spider http://localhost/ || exit 1
HealthInterval=30s
HealthStartPeriod=30s

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=60

[Install]
WantedBy=default.target
