# quadlet/dify-worker.container
[Unit]
Description=Dify Celery Worker
After=dify-api.service

[Container]
Image=docker.io/langgenius/dify-api:latest
ContainerName=dify-worker
Pod=dify-stack.pod

# Same image as API, worker mode activated via MODE=worker
EnvironmentFile=%h/.local/share/ai-models/dify/dify.env

Environment=MODE=worker
Environment=FLASK_DEBUG_MODE=false

Environment=DB_USERNAME=dify
Environment=DB_HOST=localhost
Environment=DB_PORT=5432
Environment=DB_DATABASE=dify

Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379
Environment=REDIS_USE_SSL=false
Environment=REDIS_DB=0

Environment=STORAGE_TYPE=local
Environment=STORAGE_LOCAL_PATH=/app/api/storage

Environment=VECTOR_STORE=qdrant
Environment=QDRANT_URL=http://localhost:6333
Environment=QDRANT_CLIENT_TIMEOUT=20
Environment=QDRANT_GRPC_ENABLED=false

Environment=CODE_EXECUTION_ENDPOINT=http://localhost:8194
Environment=SSRF_PROXY_HTTP_URL=http://localhost:3128
Environment=SSRF_PROXY_HTTPS_URL=http://localhost:3128
Environment=LOG_LEVEL=INFO

Volume=%h/.local/share/ai-models/dify/storage:/app/api/storage:Z

[Service]
Restart=on-failure
RestartSec=15
TimeoutStartSec=300

[Install]
WantedBy=default.target
