# quadlet/dify-postgres.container
[Unit]
Description=Dify PostgreSQL Database
After=network-online.target

[Container]
Image=docker.io/postgres:15-alpine
ContainerName=dify-postgres
Pod=dify-stack.pod

Environment=POSTGRES_USER=dify
Environment=POSTGRES_DB=dify
Environment=PGDATA=/var/lib/postgresql/data/pgdata

# DB_PASSWORD comes from the generated dify.env
# dify.env must contain POSTGRES_PASSWORD= directly — systemd does not expand
# ${VAR} from EnvironmentFile inside Environment= lines (it's a literal string).
EnvironmentFile=%h/.local/share/ai-models/dify/dify.env

Volume=%h/.local/share/ai-models/dify/postgres:/var/lib/postgresql/data:Z

HealthCmd=pg_isready -U dify
HealthInterval=30s
HealthStartPeriod=30s

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
