# quadlet/dify-web.container
[Unit]
Description=Dify Web Frontend (Next.js)
After=dify-api.service

[Container]
Image=docker.io/langgenius/dify-web:latest
ContainerName=dify-web
Pod=dify-stack.pod

# Empty → browser uses same origin; nginx routes /console/api and /api correctly
Environment=CONSOLE_API_URL=
Environment=APP_API_URL=
Environment=NEXT_TELEMETRY_DISABLED=1
# Without this Next.js binds to the pod IP only, making it unreachable via
# nginx's 127.0.0.1 upstream (all containers share the pod network namespace)
Environment=HOSTNAME=0.0.0.0

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=180

[Install]
WantedBy=default.target
