# quadlet/dify-web.container
[Unit]
Description=Dify Web Frontend (Next.js)
After=dify-api.service

[Container]
Image=docker.io/langgenius/dify-web:latest
ContainerName=dify-web
Pod=dify-stack.pod

# Empty → browser uses same origin; nginx routes /console/api and /api correctly
Environment=CONSOLE_API_URL=
Environment=APP_API_URL=
Environment=NEXT_TELEMETRY_DISABLED=1

HealthCmd=curl -sf http://localhost:3000 || exit 1
HealthInterval=30s
HealthStartPeriod=60s

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=180

[Install]
WantedBy=default.target
