# quadlet/dify-qdrant.container
[Unit]
Description=Qdrant Vector Database
After=network-online.target

[Container]
Image=docker.io/qdrant/qdrant:latest
ContainerName=dify-qdrant
Pod=dify-stack.pod

Volume=%h/.local/share/ai-models/qdrant:/qdrant/storage:Z

HealthCmd=curl -sf http://localhost:6333/readyz || exit 1
HealthInterval=30s
HealthStartPeriod=20s

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=60

[Install]
WantedBy=default.target
