# quadlet/searxng.container
[Unit]
Description=SearXNG Metasearch Engine (Web Search for AI)
After=network-online.target

[Container]
Image=docker.io/searxng/searxng:latest
ContainerName=searxng
Network=ai-stack.network
PublishPort=8888:8080

# No GPU needed — SearXNG is a lightweight web app

# Persistent config (settings.yml with limiter disabled for API use)
Volume=%h/.local/share/ai-models/searxng:/etc/searxng:Z

# Health check
HealthCmd=wget -qO- http://localhost:8080/healthz || exit 1
HealthInterval=30s
HealthStartPeriod=30s

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
